# templates/default/supervisord.conf.erb
[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port=*:9001
;username=super
;password=super

[supervisord]
logfile=<%= node["webscreenshots"]["home"] %>/var/log/supervisord.log
pidfile=<%= node["webscreenshots"]["home"] %>/var/etc/supervisord.pid
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:redis]
command=<%= node["webscreenshots"]["home"] %>/bin/redis-server
stdout_logfile=<%= node["webscreenshots"]["home"] %>/var/log/redis.log
redirect_stderr=true
priority=10

; check out https://github.com/celery/celery/tree/master/extra/supervisord
[program:celeryd]
directory=<%= node["webscreenshots"]["home"] %>
command=<%= node["webscreenshots"]["home"] %>/bin/celery worker -c 4 -E -A webscreenshots.celeryapp -l info
stdout_logfile=<%= node["webscreenshots"]["home"] %>/var/log/celeryd.log
redirect_stderr=true
startsecs=10
stopwaitsecs=300
;stopsignal=KILL
killasgroup=true
priority=20

; check out https://github.com/celery/celery/tree/master/extra/supervisord
[program:celerybeat]
directory=<%= node["webscreenshots"]["home"] %>
command=<%= node["webscreenshots"]["home"] %>/bin/celery beat -A webscreenshots.celeryapp -l info
stdout_logfile=<%= node["webscreenshots"]["home"] %>/var/log/celerybeat.log
redirect_stderr=true
startsecs=10
;stopwaitsecs=300
;stopsignal=KILL
priority=30

[program:flower]
directory=<%= node["webscreenshots"]["home"] %>
command=<%= node["webscreenshots"]["home"] %>/bin/celery flower --broker=redis://localhost:6379/0
stdout_logfile=<%= node["webscreenshots"]["home"] %>/var/log/flower.log
redirect_stderr=true
priority=40

[program:uwsgi_django]
directory=<%= node["webscreenshots"]["home"] %>
command=<%= node["webscreenshots"]["home"] %>/bin/uwsgi --ini <%= node["webscreenshots"]["home"] %>/etc/uwsgi.ini
stdout_logfile=<%= node["webscreenshots"]["home"] %>/var/log/uwsgi_django.log
redirect_stderr=true
environment=DJANGO_SETTINGS_MODULE='settings'
autostart=true
autorestart=true
priority=50
stopsignal=QUIT
