# templates/default/supervisord.conf.erb
[unix_http_server]
file=/tmp/supervisor.sock
;username=super
;password=duper

[inet_http_server]
port=*:9001
;username=super
;password=duper

[supervisord]
logfile=<%= node['webscreenshots']['venv_home'] %>/var/log/supervisord.log
pidfile=<%= node['webscreenshots']['venv_home'] %>/var/run/supervisord.pid
identifier=<%= node['hostname'] %>
logfile_maxbytes=200MB
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
;serverurl=http://localhost:9001
prompt=<%= node['hostname'] %>
history_file=<%= node['webscreenshots']['venv_home'] %>/.supervisorctl_history

[group:celery]
programs=redis,celeryd_default,celeryd_webscreenshots,celerybeat,flower
autostart=false

[program:redis]
command=/usr/local/bin/redis-server
autorestart=true
stdout_logfile=<%= node['webscreenshots']['venv_home'] %>/var/log/%(program_name)s.log
redirect_stderr=true
priority=10

; check out https://github.com/celery/celery/tree/master/extra/supervisord
;[program:celeryd_webscreenshots]
;environment=DJANGO_SETTINGS_MODULE='<%= node['webscreenshots']['django_settings_module'] %>',PYTHONPATH='.:$PYTHONPATH'
;directory=<%= node['webscreenshots']['project_root'] %>/src/webscreenshots
;command=<%= node['webscreenshots']['venv_home'] %>/bin/celery worker
;    -c 10
;    -E
;    -A webscreenshots.celeryapp
;    -l info
;    -n casperjs
;    -Q casperjs
;autorestart=true
;stdout_logfile=<%= node['webscreenshots']['venv_home'] %>/var/log/%(program_name)s.log
;redirect_stderr=true
;startsecs=5
;stopwaitsecs=60
;stopsignal=TERM
;priority=20

; check out https://github.com/celery/celery/tree/master/extra/supervisord
[program:celeryd_default]
environment=DJANGO_SETTINGS_MODULE='<%= node['webscreenshots']['django_settings_module'] %>',PYTHONPATH='.:$PYTHONPATH'
directory=<%= node['webscreenshots']['project_root'] %>/src/webscreenshots
command=<%= node['webscreenshots']['venv_home'] %>/bin/celery worker
    -c 4
    -E
    -A webscreenshots.celeryapp
    -l info
    -n default
autorestart=true
stdout_logfile=<%= node['webscreenshots']['venv_home'] %>/var/log/%(program_name)s.log
redirect_stderr=true
startsecs=5
stopwaitsecs=60
stopsignal=TERM
priority=20

; check out https://github.com/celery/celery/tree/master/extra/supervisord
[program:celerybeat]
environment=DJANGO_SETTINGS_MODULE=<%= node['webscreenshots']['django_settings_module'] %>
directory=<%= node['webscreenshots']['project_root'] %>/src/webscreenshots
; from https://github.com/Supervisor/supervisor/issues/140
command=<%= node['webscreenshots']['venv_home'] %>/bin/celery beat
    -A webscreenshots.celeryapp
    --pidfile=<%= node['webscreenshots']['venv_home'] %>/var/run/%(program_name)s.pid
    --schedule=<%= node['webscreenshots']['venv_home'] %>/var/run/celerybeat-schedule.db
    -l info
stdout_logfile=<%= node['webscreenshots']['venv_home'] %>/var/log/%(program_name)s.log
redirect_stderr=true
autorestart=true
autostart=false
startsecs=5
stopwaitsecs=10
stopsignal=TERM
priority=30

[program:flower]
directory=<%= node['webscreenshots']['venv_home'] %>
command=<%= node['webscreenshots']['venv_home'] %>/bin/celery flower --broker=redis://localhost:6379/0
stdout_logfile=<%= node['webscreenshots']['venv_home'] %>/var/log/%(program_name)s.log
autorestart=true
redirect_stderr=true
priority=40

[program:uwsgi_django]
environment=DJANGO_SETTINGS_MODULE=<%= node['webscreenshots']['django_settings_module'] %>
directory=<%= node['webscreenshots']['venv_home'] %>
command=bin/uwsgi --ini etc/uwsgi.ini
stdout_logfile=<%= node['webscreenshots']['venv_home'] %>/var/log/%(program_name)s.log
redirect_stderr=true
autorestart=true
priority=50
; this is very important!
stopsignal=INT
